# superlists/.config/supervisord.conf
[program:gunicorn]
command=gunicorn superlists.wsgi -c .config/gunicorn_cfg.py
directory=/srv/projects/superlists

[program:nginx]
command = nginx -g 'daemon off;'